# iso/airootfs/usr/lib/systemd/system/firmus-snapshot.service
[Unit]
Description=Firmus Linux Snapshot — Pre-update system snapshot
Documentation=https://github.com/firmuslinux/firmuslinux
Before=systemd-update.service

[Service]
Type=oneshot
ExecStart=/usr/bin/snapper create -d "pre-update" --description "Automatic pre-update snapshot"
ExecStartPost=/usr/bin/snapper list
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
