# iso/airootfs/usr/lib/systemd/system/firmus-snapshot.timer
[Unit]
Description=Firmus Linux Snapshot Timer — Hourly automatic snapshots

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
