# iso/airootfs/usr/lib/systemd/system/firmus-tcr-sync.timer
[Unit]
Description=Firmus Linux TCR Repository Sync Timer

[Timer]
OnCalendar=weekly
Persistent=true
RandomizedDelaySec=3600

[Install]
WantedBy=timers.target
